# deploy/topic-tracker.service
[Unit]
Description=Topic Tracker FastAPI Backend
After=network.target

[Service]
User=root
WorkingDirectory=/root/topic-tracker
ExecStart=/root/topic-tracker/venv/bin/uvicorn Backend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5
EnvironmentFile=/root/topic-tracker/Backend/.env

[Install]
WantedBy=multi-user.target
